# trt-edgellm.service
[Unit]
Description=TensorRT Edge LLM HTTP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/agi/cosmos/TensorRT-Edge-LLM
ExecStart=/home/agi/cosmos/TensorRT-Edge-LLM/build/inference_server/llm_http_server \
  --engineDir /home/agi/cosmos/TensorRT-Edge-LLM/engines/Cosmos-Reason2-2B-fp8 \
  --multimodalEngineDir /home/agi/cosmos/TensorRT-Edge-LLM/visual_engines/Cosmos-Reason2-2B-fp8 \
  --host 0.0.0.0 \
  --port 8080
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
